# docker/all-in-one/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/data/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/var/run/supervisord.pid
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; ── PostgreSQL ─────────────────────────────────────────────────────────────────
[program:postgresql]
command=gosu postgres /usr/lib/postgresql/17/bin/postgres -D /data/db/postgres
priority=10
autostart=true
autorestart=true
stopwaitsecs=60
stdout_logfile=/data/logs/postgresql.log
stderr_logfile=/data/logs/postgresql.log

; ── Mosquitto MQTT broker ──────────────────────────────────────────────────────
[program:mosquitto]
command=/usr/sbin/mosquitto -c /etc/mosquitto/conf.d/garagestack.conf
priority=20
autostart=true
autorestart=true
stdout_logfile=/data/logs/mosquitto.log
stderr_logfile=/data/logs/mosquitto.log

; ── SAIC MQTT Gateway (polls MG iSmart API, publishes to Mosquitto) ────────────
; Delayed start so Mosquitto is ready before connecting
[program:saic-gateway]
command=/bin/bash -c "sleep 10 && exec /opt/saic-venv/bin/mqtt_gateway"
priority=30
autostart=true
autorestart=true
startretries=5
stdout_logfile=/data/logs/saic-gateway.log
stderr_logfile=/data/logs/saic-gateway.log

; ── GarageStack API (.NET) ────────────────────────────────────────────────────
; Delayed start so PostgreSQL has time to accept connections
[program:api]
command=/bin/bash -c "sleep 20 && exec dotnet /app/api/GarageStack.Api.dll"
priority=40
autostart=true
autorestart=true
startretries=10
stdout_logfile=/data/logs/api-supervisor.log
stderr_logfile=/data/logs/api-supervisor.log

; ── GarageStack Worker (.NET) ─────────────────────────────────────────────────
[program:worker]
command=/bin/bash -c "sleep 25 && exec dotnet /app/worker/GarageStack.Worker.dll"
priority=50
autostart=true
autorestart=true
startretries=10
stdout_logfile=/data/logs/worker-supervisor.log
stderr_logfile=/data/logs/worker-supervisor.log

; ── Nginx (serves frontend SPA, proxies /api to the .NET API) ─────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=60
autostart=true
autorestart=true
stdout_logfile=/data/logs/nginx.log
stderr_logfile=/data/logs/nginx.log
