# docker/all-in-one/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/data/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/var/run/supervisord.pid
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; ── PostgreSQL ─────────────────────────────────────────────────────────────────
[program:postgresql]
command=gosu postgres /usr/lib/postgresql/18/bin/postgres -D /data/db/postgres
priority=10
autostart=true
autorestart=true
stopwaitsecs=60
stdout_logfile=/data/logs/postgresql.log
stderr_logfile=/data/logs/postgresql.log

; ── Mosquitto MQTT broker ──────────────────────────────────────────────────────
[program:mosquitto]
command=gosu mosquitto /usr/sbin/mosquitto -c /etc/mosquitto/conf.d/garagestack.conf
priority=20
autostart=true
autorestart=true
stdout_logfile=/data/logs/mosquitto.log
stderr_logfile=/data/logs/mosquitto.log

; ── SAIC MQTT Gateway (polls MG iSmart API, publishes to Mosquitto) ────────────
; Polls Mosquitto's port instead of a fixed sleep before connecting -- see wait-for-mosquitto.sh.
; Runs via the venv copied from the official saicismartapi Docker image.
[program:saic-gateway]
command=/bin/bash -c "/wait-for-mosquitto.sh; cd /opt/saic/app && exec gosu appuser /opt/saic/app/.venv/bin/python ./main.py"
priority=30
autostart=true
autorestart=true
startretries=5
stdout_logfile=/data/logs/saic-gateway.log
stderr_logfile=/data/logs/saic-gateway.log

; ── GarageStack API (.NET) ────────────────────────────────────────────────────
; Polls PostgreSQL's port instead of a fixed sleep before starting -- see wait-for-postgres.sh.
; directory= must stay at the DLL's own folder: ASP.NET Core resolves its content root
; (where it looks for appsettings.json) from the process's working directory, not its
; executable path, so pointing this anywhere else breaks config loading. The "logs" and
; "keys" folders app-relative paths resolve to are symlinked into /data in the Dockerfile.
[program:api]
directory=/app/api
environment=HOME="/home/appuser",DOTNET_CLI_TELEMETRY_OPTOUT="1",DOTNET_NOLOGO="1"
command=/bin/bash -c "/wait-for-postgres.sh; exec gosu appuser dotnet /app/api/GarageStack.Api.dll"
priority=40
autostart=true
autorestart=true
startretries=10
stdout_logfile=/data/logs/api-supervisor.log
stderr_logfile=/data/logs/api-supervisor.log

; ── GarageStack Worker (.NET) ─────────────────────────────────────────────────
[program:worker]
directory=/app/worker
environment=HOME="/home/appuser",DOTNET_CLI_TELEMETRY_OPTOUT="1",DOTNET_NOLOGO="1"
command=/bin/bash -c "/wait-for-postgres.sh; exec gosu appuser dotnet /app/worker/GarageStack.Worker.dll"
priority=50
autostart=true
autorestart=true
startretries=10
stdout_logfile=/data/logs/worker-supervisor.log
stderr_logfile=/data/logs/worker-supervisor.log

; ── Nginx (serves frontend SPA, proxies /api to the .NET API) ─────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=60
autostart=true
autorestart=true
stdout_logfile=/data/logs/nginx.log
stderr_logfile=/data/logs/nginx.log
